# metrics/elastos-metrics/etc/cron.d/elastos-metrics
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/2 * * * * prometheus /usr/local/bin/elastos-metrics.py
