# metrics/elastos-metrics_1.0.0-1/DEBIAN/preinst
#!/bin/sh

sysctl -w kernel.perf_event_paranoid=0
systemctl enable prometheus prometheus-node-exporter prometheus-pushgateway prometheus-alertmanager