# eth/elastos-eth/DEBIAN/postinst
#!/bin/sh

adduser --system \
        --quiet \
        --group \
        --no-create-home \
        --disabled-password \
        --home /data/elastos \
        elauser

mkdir -p /data/elastos/eth
chown -R elauser:elauser /data/elastos/eth
systemctl enable elastos-eth elastos-eth-oracle